# devel/etc/catkin/profile.d/05.catkin-test-results.sh
# generated from catkin/cmake/env-hooks/05.catkin-test-results.sh.develspace.in

export CATKIN_TEST_RESULTS_DIR="/home/sebasgm/navigation/src/nav_test/turtlebot_navigation/test_results"
export ROS_TEST_RESULTS_DIR="$CATKIN_TEST_RESULTS_DIR"
